# SimpleCDD/local_extras/mycfg/.config/autostart/tint2.desktop
[Desktop Entry]
Type=Application
Name=tint2
Exec=tint2
NoDisplay=true
