# SimpleCDD/local_extras/mycfg/.config/autostart/xscreensaver.desktop
[Desktop Entry]
Type=Application
Name=xscreensaver
Exec=xscreensaver
NoDisplay=false
