# SimpleCDD/local_extras/mycfg/etc/getty@tty1.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin da --noclear %I 38400 linux
